# 1_system_and_hardware/06_Linux_boot_Systemd/conf/greplog.timer
[Unit]
Description=run greplog.service every 30 sec

[Timer]
OnUnitActiveSec=30s
Unit=greplog.service